# base/ignore/systemd/rkhunter.timer
[Unit]
Description=Weekly Rootkit Hunter Check

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
