# base/ignore/systemd/rkhunter.service
[Unit]
Description=Rootkit Hunter Check

[Service]
Type=oneshot
ExecStart=/usr/bin/rkhunter --check --skip-keypress --quiet --cronjob
